# deploy/tradex-dashboard-unified.service
[Unit]
Description=TradeX Unified Streamlit Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tradex
ExecStart=/opt/tradex/.venv/bin/streamlit run dashboard/app_unified.py \
    --server.port=8502 \
    --server.address=0.0.0.0 \
    --server.headless=true \
    --browser.gatherUsageStats=false \
    --theme.primaryColor="#2196f3" \
    --theme.backgroundColor="#0e1117" \
    --theme.secondaryBackgroundColor="#1a1d23" \
    --theme.textColor="#fafafa"
EnvironmentFile=/opt/tradex/.env
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
